# deploy/mirage-server.service
[Unit]
Description=Mirage Server (Anti-Censorship Proxy)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Secrets live in /etc/mirage/env (mode 0600 root-only). At minimum set:
#   MIRAGE_PSK=<your-psk>
# Plus, when REALITY or an outbound relay is enabled:
#   MIRAGE_REALITY_PRIVATE_KEY=<base64-32-byte-x25519>
#   MIRAGE_OUTBOUND_UUID=<vmess-uuid>
# Keeping these off the ExecStart line stops `ps`, /proc/<pid>/cmdline
# and `systemctl cat` from leaking them.
EnvironmentFile=/etc/mirage/env
ExecStart=/opt/mirage/mirage-server \
    --domain sla.sunlawai.com \
    --web-root /opt/mirage/web \
    --cert-dir /opt/mirage/certs \
    --listen :443
Restart=always
RestartSec=3
LimitNOFILE=65535

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/opt/mirage/certs
PrivateTmp=true

[Install]
WantedBy=multi-user.target
